# update-motd.d/00-header
#!/bin/sh
cat /home/l/.config/art/kms.txt
